# projectName.ini
[program:projectname]
directory = [directoty_to_run]
command = [commant_to_run]
stdout_logfile=/var/log/supervisor/projectname.log 